# deploy/supervisor/webaudit-worker.conf
[program:webaudit-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/webaudit/current/artisan queue:work --queue=mail,default --tries=5 --backoff=60 --timeout=120 --sleep=3
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
numprocs=1
user=deploy
redirect_stderr=true
stdout_logfile=/var/www/webaudit/current/storage/logs/worker.log
